# docker/sim/rootfs/etc/supervisor/conf.d/supervisord.conf
[supervisord]
redirect_stderr=true
stopsignal=QUIT
autorestart=true
directory=/root

[group:x]
programs=xorg,wm,lxpanel,pcmanfm,x11vnc,xrandr

[program:wm]
priority=15
command=/usr/bin/openbox
user={{USER}}
environment=DISPLAY=":1",HOME="{{HOME}}",USER="{{USER}}"

[program:ibus]
priority=15
command=/usr/bin/ibus-daemon -d
user={{USER}}
environment=DISPLAY=":1",HOME="{{HOME}}",USER="{{USER}}"

[program:lxpanel]
priority=15
directory={{HOME}}
command=/usr/bin/lxpanel --profile LXDE
user={{USER}}
environment=DISPLAY=":1",HOME="{{HOME}}",USER="{{USER}}"

[program:pcmanfm]
priority=15
directory={{HOME}}
command=/usr/bin/pcmanfm --desktop --profile LXDE
user={{USER}}
stopwaitsecs=3
environment=DISPLAY=":1",HOME="{{HOME}}",USER="{{USER}}"

[program:xorg]
priority=10
command=/usr/bin/Xorg :1
stopsignal=KILL

[program:x11vnc]
priority=20
command=x11vnc -display :1 -xkb -forever -shared -repeat -capslock -noxrecord
user={{USER}}

[program:xrandr]
priority=17
command=/bin/sh /usr/local/bin/xrandr_exec.sh

[program:xcompmgr]
priority=25
command=/usr/bin/xcompmgr